# docker/roles/vagrant/templates/vagrant-mounted.conf.j2
start on vagrant-mounted

exec sudo /etc/init.d/httpd start